# potassium-image-server.service
[Unit]
Description=Potassium Local Image Server
After=network.target
Requires=network.target

[Service]
Type=simple
User=rnd2
Group=camera001
WorkingDirectory=/home/rnd2/Desktop/radar_sys/potassium-frontend
Environment=NODE_ENV=production
Environment=LOCAL_BASE_PATH=/srv/camera_uploads/camera001/192.168.1.54
ExecStart=/usr/bin/node local-image-server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=potassium-image-server

[Install]
WantedBy=multi-user.target
