# system/potassium-frontend.service
[Unit]
Description=Potassium Radar Frontend React App
After=network.target potassium-backend.service
Wants=potassium-backend.service

[Service]
Type=simple
User=rnd2
Group=rnd2
WorkingDirectory=/home/rnd2/Desktop/radar_sys/potassium-frontend
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=REACT_APP_BACKEND_URL=http://localhost:3000
Environment=GENERATE_SOURCEMAP=false
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=potassium-frontend

[Install]
WantedBy=multi-user.target
