# system/potassium-backend.service
[Unit]
Description=Potassium Radar Backend API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=rnd2
Group=rnd2
WorkingDirectory=/home/rnd2/Desktop/radar_sys/potassium-backend-
Environment=NODE_ENV=production
Environment=DB_HOST=localhost
Environment=DB_PORT=3306
Environment=DB_NAME=potassium_backend
Environment=DB_USER=root
Environment=DB_PASSWORD=RootPass2025
Environment=JWT_SECRET=your_super_secret_jwt_key_for_potassium_radar_system_2024
Environment=JWT_EXPIRES_IN=24h
Environment=PORT=3000
Environment=IMAGE_BASE_DIR=/srv/camera_uploads
Environment=FTP_HOST=192.168.1.55
Environment=FTP_PORT=21
Environment=FTP_USER=camera001
Environment=FTP_PASSWORD=RadarCamera01
Environment=UDP_LOCAL_PORT=17081
Environment=CAMERA_ID=camera001
Environment=RADAR_CAMERA_ID=RadarCamera001
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=potassium-backend

[Install]
WantedBy=multi-user.target
